# deploy/deb/openshell-gateway.service
[Unit]
Description=OpenShell Gateway
Documentation=https://github.com/NVIDIA/OpenShell
After=default.target

[Service]
Type=simple
StateDirectory=openshell/gateway
# %S resolves to $XDG_STATE_HOME for user services.
Environment=OPENSHELL_BIND_ADDRESS=127.0.0.1
Environment=OPENSHELL_SERVER_PORT=17670
Environment=OPENSHELL_TLS_CERT=%S/openshell/tls/server/tls.crt
Environment=OPENSHELL_TLS_KEY=%S/openshell/tls/server/tls.key
Environment=OPENSHELL_TLS_CLIENT_CA=%S/openshell/tls/ca.crt
Environment=OPENSHELL_DB_URL=sqlite:%S/openshell/gateway/openshell.db
Environment=OPENSHELL_GRPC_ENDPOINT=https://127.0.0.1:17670
Environment=OPENSHELL_SSH_GATEWAY_HOST=127.0.0.1
Environment=OPENSHELL_SSH_GATEWAY_PORT=17670
Environment=OPENSHELL_VM_DRIVER_STATE_DIR=%S/openshell/vm-driver
Environment=OPENSHELL_VM_TLS_CA=%S/openshell/tls/ca.crt
Environment=OPENSHELL_VM_TLS_CERT=%S/openshell/tls/client/tls.crt
Environment=OPENSHELL_VM_TLS_KEY=%S/openshell/tls/client/tls.key
Environment=OPENSHELL_DOCKER_TLS_CA=%S/openshell/tls/ca.crt
Environment=OPENSHELL_DOCKER_TLS_CERT=%S/openshell/tls/client/tls.crt
Environment=OPENSHELL_DOCKER_TLS_KEY=%S/openshell/tls/client/tls.key
Environment=OPENSHELL_PODMAN_TLS_CA=%S/openshell/tls/ca.crt
Environment=OPENSHELL_PODMAN_TLS_CERT=%S/openshell/tls/client/tls.crt
Environment=OPENSHELL_PODMAN_TLS_KEY=%S/openshell/tls/client/tls.key
EnvironmentFile=-%h/.config/openshell/gateway.env
ExecStartPre=/usr/bin/openshell-gateway generate-certs --output-dir %S/openshell/tls --server-san host.openshell.internal
ExecStart=/usr/bin/openshell-gateway
Restart=on-failure
RestartSec=5s
PrivateTmp=true
UMask=0077

[Install]
WantedBy=default.target
